# conf/app.conf
[program:gunicorn_app]
command=/opt/app/exec.sh
